# arb-bot/deploy/gold-v2-monitor.service
[Unit]
Description=Gold V2 arbitrage safety monitor
After=xau-arb-bot.service network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/root/perp-arb-bot/arb-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u scripts/monitor_gold_v2.py --cycles 3 --interval 3 --max-minutes 1440 --log-file /root/perp-arb-bot/arb-bot/data/gold_v2_monitor_systemd.log
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
