# arb-bot/deploy/xau-arb-bot.service
[Unit]
Description=MT4 XAUUSD Binance arbitrage executor
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/root/perp-arb-bot/arb-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8011 --no-access-log
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
